# supervisor/smart-exporter.ini
[program:smart-exporter]
command=/usr/local/bin/smart-exporter
user=root
autorestart=true
redirect_stderr=true
stdout_logfile_maxbytes=20MB
stdout_logfile=/var/log/supervisord/smart-exporter.log
